# supervisord.conf
[supervisord]
nodaemon=true
environment=SPARK_NO_DAEMONIZE="1"

[program:spark]
command=/opt/spark/sbin/start-slave.sh spark://spark-master:7077
